# robotik_chess_stack/agent/systemd/agent-api.service
[Unit]
Description=Agent API Service (Robotik-Schach)
After=network.target

[Service]
Type=simple
User=agent
WorkingDirectory=/opt/agent
Environment=PYTHONUNBUFFERED=1
Environment=AGENT_CONFIG=/etc/agent/agent.yaml
Environment=PYTHONPATH=/opt/agent/src
ExecStart=/opt/agent/venv/bin/uvicorn agent_api.app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
